# pi-dashboard.service
[Unit]
Description=Cozy pixel-art Pi dashboard backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/dashboard/desktopPI
# Secrets and deployment switches live here (CANVAS_TOKEN, CALENDAR_ICS_URL, BIND).
# The leading "-" makes it optional so the unit still starts without the file.
EnvironmentFile=-/home/pi/dashboard/desktopPI/.env
# Needed by wlopm, the Wayland fallback for display on/off.
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/usr/bin/python3 /home/pi/dashboard/desktopPI/server.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
